# infra/promtail/promtail.service
# User-mode systemd unit — install to ~/.config/systemd/user/promtail.service
# Manages as: systemctl --user enable|start|restart|status promtail
[Unit]
Description=Promtail log shipper → Grafana Cloud Loki
After=network-online.target

[Service]
EnvironmentFile=%h/.config/promtail/secrets
ExecStart=%h/bin/promtail -config.file=%h/.config/promtail/config.yml -config.expand-env=true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
