# infra/promtail/promtail.service
[Unit]
Description=Promtail log shipper → Grafana Cloud Loki
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/promtail/secrets
ExecStart=/usr/local/bin/promtail -config.file=/etc/promtail/config.yml -config.expand-env=true
Restart=on-failure
RestartSec=5s
User=promtail
Group=promtail

[Install]
WantedBy=multi-user.target
